# huita8/API/task.cron
1 * * * * wget -q -O - http://10.103.1.202/~huitahuit/huita8/API/?API=hello